# rc-services/openbox
#!/sbin/openrc-run

name="openbox"
description="Openbox window manager"

command="/usr/bin/openbox"
command_args="--startup /etc/xdg/openbox/autostart"
command_background="yes"
pidfile="/run/${RC_SVCNAME}.pid"

command_user="lms:lms"

export DISPLAY=":99"

log_dir="/var/log/openbox"
output_log="${log_dir}/openbox.log"
error_log="${log_dir}/openbox.err"

depend() {
    need xvfb
    after xvfb
}

start_pre() {
    checkpath --directory --mode 0755 /run
    checkpath --directory --mode 0755 --owner lms:lms "${log_dir}"
}

start_post() {
    ebegin "openbox started on ${DISPLAY}"
    eend 0
}

stop_post() {
    ebegin "openbox stopped"
    eend 0
}
