# rc-services/xterm
#!/sbin/openrc-run

name="xterm"
command="/usr/bin/xterm"
command_user="lms"
command_background="yes"
export DISPLAY=":0"
export HOME="/home/lms"

depend() {
    need openbox
}
